# packages/xinstall/xinstall-0.52.2.tar.gz/xinstall-0.52.2/xinstall/data/linux/autostart/autolock.desktop
[Desktop Entry]
Comment=Automatically locks screen after login
Encoding=UTF-8
Exec=sleep 5s && xscreensaver-command -lock
Hidden=true
Name=Auto Lock
StartupNotify=false
Terminal=false
Type=Application
Version=0.9.4
