# packages/xinstall/xinstall-0.52.2-py3-none-any.whl/xinstall/data/linux/autostart/pia_manager.desktop
[Desktop Entry]
Type=Application
Name=Private Internet Access
Exec=/home/dclong/.pia_manager/pia_manager/run.sh --startup
X-GNOME-Autostart-enabled=true
Icon=/home/dclong/.pia_manager/pia_tray.32/Resources/default_app_logo.png
