# packages/xinstall/xinstall-0.52.2-py3-none-any.whl/xinstall/data/linux/autostart/setxkbmap.desktop
[Desktop Entry]
Encoding=UTF-8
Version=0.9.4
Type=Application
Name=SetXKBMap
Comment=Map the CapsLock key to the Escape key
Exec=setxkbmap -option -option caps:escape
OnlyShowIn=XFCE;
StartupNotify=false
Terminal=false
Hidden=false

