# packages/xinstall/xinstall-0.52.2-py3-none-any.whl/xinstall/data/linux/autostart/terminator.desktop
[Desktop Entry]
Type=Application
Name=Terminator
Comment=Run Terminator
Exec=terminator -l s2 
StartupNotify=false
Terminal=false
Hidden=false

